# salt-state-files/history-date-time/profile.d/history.sh
#!/bin/sh
export HISTTIMEFORMAT='%F %t'
